# data/plasma-xdg-desktop-portal-kde.service.in
[Unit]
Description=Xdg Desktop Portal For KDE/SonicDE (X11)
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/xdg-desktop-portal-kde
BusName=org.freedesktop.impl.portal.desktop.kde
Slice=session.slice
# Auto restart is off as we are DBus activated
Restart=no
